# api/autocomplete/autocomplete_update.crontab
MAILTO=troy@beercrush.com

# m h dom mon dow user    command
0 6 * * * www-data /usr/local/beercrush/bin/autocomplete_update -C /etc/BeerCrush/webapp.conf
