# conf/supervisord/autocomplete_watch.conf
[program:autocomplete_watch]
command=/usr/local/beercrush/bin/autocomplete_watch -C /etc/BeerCrush/webapp.conf
autostart=true
redirect_stderr=true
user=www-data
