# conf/supervisord/couchdb-proxy.conf
[program:spread-couchdb-proxy]
command=/usr/local/beercrush/spread-couchdb-proxy
autostart=true
redirect_stderr=true
user=www-data
