# conf/supervisord/memcacheq.conf
[program:memcacheq]
command=memcacheq -u troy -r -H /var/lib/memcacheq -N -v -L 1024 -B 1024
autostart=true
redirect_stderr=true

