# conf/supervisord/php-cgi.conf
[program:spread-php-cgi]
command=/usr/local/beercrush/spread-php-cgi
autostart=true
redirect_stderr=true
user=www-data
