# conf/supervisord/queueproc.conf
[program:queueproc]
command=/usr/local/beercrush/queueproc
environment=OAKConfig=/etc/BeerCrush/webapp.conf
autostart=true
redirect_stderr=true
user=www-data
