# conf/supervisord/watch_changes.conf
[program:watch_changes]
command=/usr/local/beercrush/watch_changes
environment=OAKConfig=/etc/BeerCrush/webapp.conf
autostart=true
redirect_stderr=true
user=www-data
