# conf/supervisord/web.conf
[program:spread-web]
command=/usr/local/beercrush/spread-web
autostart=true
redirect_stderr=true
user=www-data
